# etc/s6-overlay/s6-rc.d/proxmox-backup/run
#!/bin/bash

# remove old lock files
rm -f /etc/proxmox-backup/.*.lck /etc/proxmox-backup/*.lock

exec /usr/lib/*/proxmox-backup/proxmox-backup-api
