# etc/s6-overlay/s6-rc.d/proxmox-backup-proxy/run
#!/bin/bash

exec s6-setuidgid backup \
    /usr/lib/*/proxmox-backup/proxmox-backup-proxy
